# init/indieauth.service
[Unit]
After=network.target
After=syslog.target
Description=Personal IndieAuth server instance
Documentation=https://indieauth.net/source/

[Service]
EnvironmentFile=/etc/indieauth/env
ExecStart=/usr/local/bin/indieauth
Group=indieweb
Restart=always
RestartSec=2s
Type=simple
User=indieweb
WorkingDirectory=/var/lib/indieauth/

[Install]
Alias=indieauth
WantedBy=multi-user.target
